# src/arch-repo-grug.timer
[Unit]
Description=arch-repo-grug update timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
